# integrations/openclaw/chainwatch-intercept.service
[Unit]
Description=Chainwatch LLM Intercept Proxy
Documentation=https://github.com/ppiankov/chainwatch
After=network-online.target
Wants=network-online.target
Before=openclaw-gateway.service

[Service]
Type=simple
ExecStart=/usr/local/bin/chainwatch intercept \
  --port 9999 \
  --upstream https://api.anthropic.com \
  --profile clawbot \
  --audit-log /var/log/chainwatch/intercept-audit.jsonl
Restart=always
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/chainwatch
ReadOnlyPaths=/root/.chainwatch
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true

MemoryMax=256M
TasksMax=50

[Install]
WantedBy=multi-user.target
